# ubuntu/base-apps/rootfs/etc/services.d/mongodb/finish
#!/usr/bin/with-contenv bash

MONGODB_ROOTDIR="$BITNAMI_PREFIX/mongodb"

# Remove pidfile
rm -f $MONGODB_ROOTDIR/tmp/mongodb.pid
